# scripts/systemd/oif-frontend.service
[Unit]
Description=Eden Portal Frontend (Next.js)
After=oif-frontend-api.service
PartOf=oif.target

[Service]
Type=simple
WorkingDirectory=/root/eden-apps
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=pnpm dev:eden
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=oif.target
