# scripts/systemd/oif-aggregator.service
[Unit]
Description=OIF Aggregator
After=network.target
PartOf=oif.target

[Service]
Type=simple
WorkingDirectory=/root/solver-cli/oif/oif-aggregator
EnvironmentFile=/root/solver-cli/.env
Environment=RUST_LOG=info
ExecStart=/root/solver-cli/oif/oif-aggregator/target/release/oif-aggregator
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=oif.target
