# scripts/systemd/oif-frontend-api.service
[Unit]
Description=Eden Portal Backend API
After=network.target oif-aggregator.service
PartOf=oif.target

[Service]
Type=simple
WorkingDirectory=/root/eden-apps/apps/eden-portal
EnvironmentFile=/root/eden-apps/apps/eden-portal/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=tsx server/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=oif.target
