# scripts/systemd/oif-aggregator.service
[Unit]
Description=OIF Aggregator
After=network.target
PartOf=oif.target

[Service]
Type=simple
WorkingDirectory=/root/solver-cli
EnvironmentFile=/root/solver-cli/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=RUST_LOG=info
ExecStart=/root/solver-cli/.aggregator/bin/oif-aggregator
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=oif.target
