# scripts/systemd/oif-frontend.service
[Unit]
Description=OIF Frontend (Vite)
After=oif-frontend-api.service
PartOf=oif.target

[Service]
Type=simple
WorkingDirectory=/root/solver-cli/frontend
ExecStartPre=npm install --silent
ExecStart=npx vite --host
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=oif.target
