# scripts/systemd/oif-solver.service
[Unit]
Description=OIF Solver
After=network.target oif-aggregator.service
PartOf=oif.target

[Service]
Type=simple
WorkingDirectory=/root/solver-cli
EnvironmentFile=/root/solver-cli/.env
ExecStart=/root/solver-cli/target/release/solver-cli solver start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=oif.target
